# maxmind/crontab
1 * * * * * * /app/maxmind/geoipupdate -f /app/maxmind/GeoIP.conf --verbose